# autoplex.service
[Unit]
Description=Completed torrent auto-extracter.

[Service]
ExecStart=/usr/local/bin/autoplex

[Install]
WantedBy=multi-user.target
